# supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=python fastapi_app.py
autostart=true
autorestart=true

[program:flask]
command=python flask_app.py
autostart=true
autorestart=true